# todos/supervisor.conf
[supervisord]
user = root

[program:uvicorn]
user = root
command = uvicorn app:app --host 0.0.0.0 --port 80
